# config/systemd/audio-sync.service
[Unit]
Description=Pi Audio Sync Service
After=network.target pipewire.service pipewire-pulse.service
Requires=pipewire.service pipewire-pulse.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Environment=PYTHONPATH=/opt/pi-audio-sync
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PIPEWIRE_RUNTIME_DIR=/run/user/1000/pipewire
ExecStartPre=/bin/sh -c 'until pw-cli list > /dev/null 2>&1; do sleep 1; done'
ExecStart=/usr/bin/python3 -u -m src.main
WorkingDirectory=/opt/pi-audio-sync
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
User=%u

[Install]
WantedBy=default.target 