# config/systemd/audio-sync.service
[Unit]
Description=Pi Audio Sync Service
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Environment=PYTHONPATH=/opt/pi-audio-sync
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/python3 -u -m src.main
WorkingDirectory=/opt/pi-audio-sync
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
User=%u

[Install]
WantedBy=default.target 