# config/systemd/audio-sync.service
[Unit]
Description=Pi Audio Sync Service
PartOf=pipewire.service
After=pipewire.service pipewire-pulse.service

[Service]
Type=simple
WorkingDirectory=/opt/pi-audio-sync
Environment=PYTHONPATH=/opt/pi-audio-sync
ExecStart=/usr/bin/python3 -m src.main
Restart=always
RestartSec=5

[Install]
WantedBy=pipewire.service 