# config/systemd/audio-sync.service
[Unit]
Description=Pi Audio Sync Service
After=network.target

[Service]
Type=simple
User=pi
Environment=PYTHONPATH=/opt/pi-audio-sync
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
ExecStartPre=/bin/sh -c 'until systemctl --user -M pi@ is-active pipewire.service; do sleep 1; done'
ExecStart=/usr/bin/python3 -u -m src.main
WorkingDirectory=/opt/pi-audio-sync
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target 