# scripts/audio-sync.service
[Unit]
Description=Pi Audio Sync Service
After=network.target pipewire.service bluetooth.service
Wants=pipewire.service bluetooth.service

[Service]
Type=simple
User=vo603-rpi-1
Group=vo603-rpi-1
WorkingDirectory=/opt/pi-audio-sync
Environment=PYTHONPATH=/opt/pi-audio-sync
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket
ExecStartPre=/usr/bin/sudo /usr/sbin/usermod -a -G bluetooth,audio vo603-rpi-1
ExecStart=/usr/bin/python3 -m src.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 